# dist/systemd/fedora-coreos-pinger.service
[Unit]
Description=Telemetry service for Fedora CoreOS
Documentation=https://github.com/coreos/fedora-coreos-pinger
Before=systemd-user-sessions.service

[Service]
DynamicUser=yes
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/fedora-coreos-pinger

[Install]
WantedBy=multi-user.target
